# daemon/client/tunnel-client.service
[Unit]
Description=Air-gapped Tunnel Client
Documentation=https://github.com/idp/tunnel
After=network.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=tunnel
Group=tunnel
WorkingDirectory=/opt/tunnel-client
ExecStart=/bin/bash -c '/opt/tunnel-client/bin/tunnel-client-linux -server="$TUNNEL_SERVER_URL" -token="$TUNNEL_TOKEN" -forward="$TUNNEL_FORWARD" -id="$TUNNEL_CLIENT_ID" -improved=true -metrics=true  -skip-verify="$TUNNEL_SKIP_VERIFY"'
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
TimeoutStopSec=30
KillMode=process
KillSignal=SIGTERM

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/tunnel-client/logs
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

# Environment file
EnvironmentFile=-/etc/tunnel-client/config

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tunnel-client

[Install]
WantedBy=multi-user.target


