# daemon/server/tunnel-server.service
[Unit]
Description=Air-gapped Tunnel Server
Documentation=https://github.com/s4l10u/tunnel
After=network.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=tunnel-server
Group=tunnel-server
WorkingDirectory=/opt/tunnel-server
ExecStart=/opt/tunnel-server/bin/tunnel-server-linux -config=/etc/tunnel-server/config.yaml
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
TimeoutStopSec=30
KillMode=process
KillSignal=SIGTERM

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/tunnel-server/logs
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

# Configuration file (YAML format - recommended)
# Falls back to environment variables if YAML config not found
# You can override config path with: TUNNEL_CONFIG=/path/to/config.yaml
Environment=TUNNEL_CONFIG=/etc/tunnel-server/config.yaml
EnvironmentFile=-/etc/tunnel-server/config

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tunnel-server

[Install]
WantedBy=multi-user.target